# blink.service
[Unit]
Description=Vamos piscar leds.

[Service]
Type=simple
ExecStart=